# chat_system.service
[Unit]
Description=Chat System Gunicorn Application
After=network.target
Wants=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=chat_system
PIDFile=/var/run/chat_system/gunicorn.pid
WorkingDirectory=/var/www/chat_system
Environment=PATH=/var/www/chat_system/venv/bin
Environment=FLASK_ENV=production
Environment=FLASK_CONFIG=production
ExecStart=/var/www/chat_system/venv/bin/gunicorn --config /var/www/chat_system/gunicorn_config.py flask.wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/chat_system/instance
ReadWritePaths=/var/log/chat_system
ReadWritePaths=/var/run/chat_system
NoNewPrivileges=true
LimitNOFILE=65536

# 重启策略
Restart=always
RestartSec=10
StartLimitIntervalSec=600
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
